# conf/conf.d/celery.conf
[program:celery]
command=celery worker -l DEBUG
user=celery
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

