# conf/conf.d/rabbitmq.conf
[program:rabbitmq]
command=rabbitmq-server
priority=1

